# services/consolekit-log/run
#!/bin/execlineb -P

s6-setuidgid 19:19
exec -c
s6-log n3 T s100000 /var/log/aux-serv/consolekit
